# files/etc/profile.d/99-custom-banner.sh
#!/bin/sh
[ -x /etc/banner.sh ] && /etc/banner.sh
